# launchd/com.example.garmincoach.dailysync.plist
<?xml version="1.0" encoding="UTF-8"?>
<!-- macOS LaunchAgent. Replace YOUR_USERNAME and the path, then:
     cp this ~/Library/LaunchAgents/
     launchctl bootstrap gui/$(id -u) ~/Library/LaunchAgents/com.example.garmincoach.dailysync.plist
     Runs daily at 07:30 (and catches up on wake if the Mac was asleep). -->
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.example.garmincoach.dailysync</string>
  <key>ProgramArguments</key>
  <array>
    <string>/Users/YOUR_USERNAME/garmin-ai-coach/venv/bin/python</string>
    <string>/Users/YOUR_USERNAME/garmin-ai-coach/scripts/daily_sync.py</string>
  </array>
  <key>WorkingDirectory</key>
  <string>/Users/YOUR_USERNAME/garmin-ai-coach</string>
  <key>StartCalendarInterval</key>
  <dict>
    <key>Hour</key><integer>7</integer>
    <key>Minute</key><integer>30</integer>
  </dict>
  <key>StandardOutPath</key>
  <string>/Users/YOUR_USERNAME/garmin-ai-coach/data/sync.log</string>
  <key>StandardErrorPath</key>
  <string>/Users/YOUR_USERNAME/garmin-ai-coach/data/sync.log</string>
  <key>RunAtLoad</key>
  <false/>
</dict>
</plist>
